# scripts/tesserae-pi-png-client.service.in
[Unit]
Description=Tesserae Pi PNG client (Pimoroni e-ink panel)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=@USER@
Group=@USER@
SupplementaryGroups=gpio spi
ExecStart=/usr/local/bin/tesserae-pi-png-client
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
